# cat-monitor.service
[Unit]
Description=Cat Monitor Motion Recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi5

# Adjust this path if you clone the repo somewhere else
WorkingDirectory=/home/pi5/fair-feeder

# If you use a virtual environment, use that python executable instead
# e.g., ExecStart=/home/pi5/fair-feeder/venv/bin/python motion_recorder.py
ExecStart=/usr/bin/python3 /home/pi5/fair-feeder/motion_recorder.py

# Auto-restart if the script crashes
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
